# deploy/supervisor/reviveguard.conf
[program:reviveguard-horizon]
process_name=%(program_name)s
command=php /var/www/reviveguard/artisan horizon
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
; Give Horizon 1 hour to finish processing current jobs before force kill

[program:reviveguard-scheduler]
process_name=%(program_name)s
; Laravel scheduler — runs every minute as a Supervisor process
; This is an alternative to system cron for hosted environments
command=php /var/www/reviveguard/artisan schedule:work
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
